# scripts/pi-print-ws/cloudflared-pi-print-ws.service
[Unit]
Description=Cloudflare quick tunnel -> pi-print-ws (8765)
After=network-online.target pi-print-ws.service
Wants=network-online.target

[Service]
Type=simple
# Copie cloudflared-quick-tunnel-sync.sh para ~/pi-print-ws/ e chmod +x. .env: SUPABASE_URL, SUPABASE_ANON_KEY, PI_TUNNEL_SYNC_SECRET (= tunnel_sync_secret no banco).
ExecStart=/usr/bin/bash /home/kim/pi-print-ws/cloudflared-quick-tunnel-sync.sh
Restart=always
RestartSec=10
User=kim
Group=kim

[Install]
WantedBy=multi-user.target
